# pi/services/kaetram-server.service
[Unit]
Description=Kaetram Game Server
After=network.target

[Service]
Type=simple
User=rjudge
WorkingDirectory=/home/rjudge/FitRPG/games/kaetram
Environment=PATH=/home/rjudge/.nvm/versions/node/v20.20.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/rjudge/.nvm/versions/node/v20.20.0/bin/node -r /home/rjudge/.nvm/versions/node/v20.20.0/lib/node_modules/tsx/dist/cli.mjs packages/server/src/main.ts
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
